# www/systemd/restart_astra.service
[Unit]
Description=Restart Astra web stack (nginx, python services, VNC/noVNC)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/var/www/restart_astra.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
