# astra-monitor/installer/configs/systemd/live-server.service
[Unit]
Description=Astra Monitor - Live Server (Flask)
After=network.target nginx-custom.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/live-server
Environment="PATH=/var/www/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin"
Environment="FLASK_APP=server.py"
ExecStart=/var/www/.venv/bin/python3 /var/www/live-server/server.py
Restart=always
RestartSec=10
StandardOutput=append:/var/www/logs/live-server.log
StandardError=append:/var/www/logs/live-server-error.log

[Install]
WantedBy=multi-user.target