# rtmp-server/debian/rtmp-server.prerm
#!/bin/sh
# Останавливает только то, что принадлежит самому RTMP-server
# (rtmp-server-gui.service — автозапуск GUI). Юниты сайта
# (nginx-rtmp/live-server/reboot-server) НИКОГДА не трогаются здесь —
# они принадлежат сайту, а не этому пакету, и должны пережить его удаление.
set -e

systemctl stop rtmp-server-gui.service >/dev/null 2>&1 || true
systemctl disable rtmp-server-gui.service >/dev/null 2>&1 || true

exit 0
